# files/ansible_player.cron
0,15,30,45 * * * * /usr/local/bin/ansible_player
