# files/ansible_player.cron
0,15,30,45 * * * * root /usr/local/bin/ansible_player >> /var/log/ansible_player.log
